# plugins/statusline-usage-updater/launchd/com.joelpt.statusline-usage-updater.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.joelpt.statusline-usage-updater</string>

    <!-- Absolute paths: LaunchAgents do not inherit shell env. -->
    <key>ProgramArguments</key>
    <array>
        <string>@@PYTHON3@@</string>
        <string>@@PLUGIN_ROOT@@/bin/update_coefficient.py</string>
    </array>

    <!-- Run once daily at 03:17 local time. Off-peak by every interpretation.
         StartCalendarInterval omits Day/Weekday so it fires daily. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>17</integer>
    </dict>

    <!-- If the Mac was asleep at 03:17, run as soon as it wakes. -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Keep tidy stderr/stdout in case launchctl logs help diagnosis. -->
    <key>StandardOutPath</key>
    <string>@@STATE_DIR@@/log/launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>@@STATE_DIR@@/log/launchd.stderr.log</string>

    <!-- /usr/bin first so /usr/bin/python3 is reachable; expand as needed. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/local/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>@@HOME@@</string>
    </dict>

    <!-- Nice the process; this is housekeeping, not a deadline. -->
    <key>Nice</key>
    <integer>5</integer>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
